# home/.config/openbox/supervisord.d/docky.ini
[program:docky]
command=docky
autostart=false
