# home/.config/openbox/supervisord.d/dzen2-time.ini
[program:time]
command=%(here)s/scripts/time.sh
autostart=false
