# home/.config/openbox/supervisord.d/synapse.ini
[program:synapse]
command=synapse --startup
