# install/files/zerotier.service
[Unit]
Description=zerotier

[Service]
ExecStart=/home/aaron/.nix-profile/bin/zerotier-one

[Install]
WantedBy=multi-user.target
